# package/etc/cron.d/overrustle
5 0 1 * * overrustle overrustle-compress-logs.sh