# package/etc/systemd/system/orl-server.service
[Unit]
Description=Overrustlelogs Server-Part
After=network.target

[Service]
User=overrustlelogs
Restart=always
EnvironmentFile=/etc/overrustlelogs/systemd.conf
ExecStart=/usr/bin/orl-server $DAEMON_OPTS

[Install]
WantedBy=multi-user.target
